# cron/cron.hourly
#!/bin/bash
perl REPLACELBPHTMLAUTHDIR/index.cgi cron=true &> /dev/null
